# packer/swarmer/src/swarmer.path
[Unit]
Description=Swarmer Configuration Path

[Path]
PathExists=/etc/swarmer/swarmer.conf
Unit=swarmer.service

[Install]
WantedBy=default.target
