# packer/consul/scripts/consul.path
[Unit]
Description=Consul Configuration

[Path]
DirectoryNotEmpty=/etc/consul
MakeDirectory=yes
