# packer/consul/scripts/docker-configurator.path
[Unit]
Description=Docker Engine Configuration Path
[Path]
DirectoryNotEmpty=/etc/docker.conf.d
MakeDirectory=yes
