# packer/consul/scripts/docker.path
[Unit]
Description=Docker Engine Configuration
[Path]
DirectoryNotEmpty=/etc/systemd/system/docker.service.d
