# packer/consul/scripts/user-config-ovfenv.path
[Unit]
Description=Hack to bypass NoSuchFile error
[Path]
DirectoryNotEmpty=/etc/ovfenvpath
MakeDirectory=yes
